# jobs/org.aesadde.productiveTrue.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>org.aesadde.productiveTrue</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/aesadde/scripts/productive.sh</string>
      <string>-true</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/tmp/productiveTrue.err</string>

    <key>StandardOutPath</key>
    <string>/tmp/productiveTrue.out</string>

    <!-- Run every day at 7 am-->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>7</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </dict>
</plist>

